# cookbooks/fimk/templates/default/fimk.conf.erb
description "FIMK server"

start on filesystem or runlevel [2345]
stop on runlevel [!2345]

setuid fim

respawn

chdir /home/fim/fim

# this application is a "resource hog"
# oom score 1000

exec java -cp fim.jar:lib/*:conf nxt.Nxt

<% if node[:ssmtp] then %>
post-start script
  TEMP_FILE=/tmp/temp_fimk_mail_<%= @to %>.txt
  rm -f $TEMP_FILE; touch $TEMP_FILE

  echo "To: <%= @to %>" >> $TEMP_FILE
  echo "From: <%= @from %>" >> $TEMP_FILE
  echo "Subject: FIMK <%= @from %>\n" >> $TEMP_FILE
  echo "STARTED" >> $TEMP_FILE
  echo <%= "http://#{node['ipaddress']}:6886/nxt?requestType=getState" %> >> $TEMP_FILE
  echo <%= "http://#{node['ipaddress']}:7886/nxt?requestType=getState" %> >> $TEMP_FILE

  ssmtp <%= @to %> < $TEMP_FILE
end script

post-stop script
  while pgrep -f 'fim.jar' > /dev/null; do sleep 1; done

  TEMP_FILE=/tmp/temp_fimk_mail_<%= @to %>.txt
  LOG_FILE=<%= @log_file %>
  rm -f $TEMP_FILE; touch $TEMP_FILE

  echo "To: <%= @to %>" >> $TEMP_FILE
  echo "From: <%= @from %>" >> $TEMP_FILE
  echo "Subject: FIMK <%= @from %>\n" >> $TEMP_FILE
  echo "STOPPED" >> $TEMP_FILE
  echo <%= "http://#{node['ipaddress']}:6886/nxt?requestType=getState" %> >> $TEMP_FILE
  echo <%= "http://#{node['ipaddress']}:7886/nxt?requestType=getState" %> >> $TEMP_FILE

  ssmtp <%= @to %> < $TEMP_FILE
  sleep 10
end script
<% else %>
post-stop script
  while pgrep -f 'fim.jar' > /dev/null; do sleep 1; done
  sleep 10
end script
<% end %>